# deploy/systemd/zenic-agents.service
[Unit]
Description=Zenic-Agents - Plataforma de Asistencia Empresarial Inteligente
Documentation=https://github.com/zenic-agents/zenic
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify
NotifyAccess=all

# User and group
User=zenic
Group=zenic
WorkingDirectory=/opt/zenic-agents

# Environment
EnvironmentFile=/opt/zenic-agents/.env

# Start command — Gunicorn + Uvicorn workers for production
ExecStart=/opt/zenic-agents/venv/bin/gunicorn \
    src.server.fastapi_app:create_app_from_env \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 0.0.0.0:5000 \
    --timeout 120 \
    --graceful-timeout 30 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --access-logfile /var/log/zenic/access.log \
    --error-logfile /var/log/zenic/error.log \
    --log-level info \
    --factory

# Graceful stop
ExecStop=/bin/kill -TERM $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true

# Read-only paths
ReadWritePaths=/home/zenic/.zenic-agents /var/log/zenic
ReadOnlyPaths=/opt/zenic-agents

# Resource limits
LimitNOFILE=65536
MemoryMax=2G

[Install]
WantedBy=multi-user.target
